# static/run-vnc-setup.service
[Unit]
Description=Run VNC setup script on first boot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-vnc.sh
User=ubuntu
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
