# templates/gnome-session-xvfb.service
[Unit]
Description=X11 Screen Mirroring from :0 to :1
After=xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=$USERN
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/$USERN/.Xauthority
ExecStart=/usr/bin/x11vnc -display :0 -viewonly -shared -forever -loop -noxdamage -repeat -clip 1920x1080+0+0 -rfbport 5901 -desktop "Mirror"
ExecStartPost=/bin/bash -c 'sleep 2 && DISPLAY=:1 XAUTHORITY=/home/$USERN/.Xauthority /usr/bin/vncviewer -fullscreen -viewonly localhost:5901'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
